# debian/debfiles/postinst
#!/bin/bash

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    systemctl daemon-reload
    systemctl enable wire-pod
    # Print web URL before starting service (same order as legacy debs; avoids set -e / stderr hiding issues).
    JUST_OUT=""
    if [ -f /usr/bin/wire-pod ]; then
        JUST_OUT=$(/usr/bin/wire-pod -justip 2>&1) || true
        printf '%s\n' "$JUST_OUT"
    else
        echo "ERROR: /usr/bin/wire-pod is missing — this .deb was built without a successful compile. Rebuild with a working Go (see build log)." >&2
    fi
    systemctl start wire-pod || true
    if systemctl is-active --quiet wire-pod 2>/dev/null; then
        echo "wire-pod is installed and the service is running."
    else
        echo "wire-pod is installed but the service is NOT running. Check: journalctl -u wire-pod -e --no-pager" >&2
    fi
    # Fallback if binary printed nothing useful (wrong arch, crash before Println, etc.)
    if ! printf '%s' "$JUST_OUT" | grep -qF 'http://'; then
        WEBP=8080
        if [ -r /etc/wire-pod/config.ini ]; then
            _wp=$(grep -E '^[[:space:]]*web_port[[:space:]]*=' /etc/wire-pod/config.ini | head -n1 | sed 's/^[^=]*=[[:space:]]*//' | tr -d '\r')
            [ -n "$_wp" ] && WEBP="$_wp"
        fi
        IP=$(hostname -I 2>/dev/null | awk '{print $1}')
        echo "WirePod configuration page: http://${IP:-127.0.0.1}:$WEBP"
    fi
    if [ ! -r /etc/wire-pod/webroot/index.html ]; then
        echo "WARNING: /etc/wire-pod/webroot/index.html missing — web UI will be empty. Reinstall or copy webroot from package source." >&2
    fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
